# etc/systemd/system/gdm.service.d/override.conf.tmpl
[Service]
ExecStartPre=/bin/cp {{ .chezmoi.homeDir }}/.config/monitors.xml /var/lib/gdm/seat0/config/monitors.xml
